# services/fan.service
[Unit]
Description=Fan Lights
After=multi-user.target

[Service]
Type=idle
Restart=always
RestartSec=1
Environment=SERVICE_ROOT=
StartLimitInterval=0
StartLimitBurst=120
ExecStart=/usr/bin/env python3 -u ${SERVICE_ROOT}/app/fan.py -q
ExecStopPost=-/usr/bin/env python3 -u ${SERVICE_ROOT}/app/fan.py --off -q
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
